# PG-Strom/supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
command=./opt/pg-strom/start_postgres.sh



