# PG-Strom/supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
command=/etc/init.d/postgresql restart



